# scripts/com.yoongi.blog-post-am.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yoongi.blog-post-am</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source "/Users/yoongi/Desktop/Study/Github/innerpyg.github.io/scripts/blog-post-today.sh" am</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yoongi/innerpyg.github.io/logs/launchagent-am-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yoongi/innerpyg.github.io/logs/launchagent-am-err.log</string>
</dict>
</plist>
